# hosts/fs.secure/netsuspend.conf
# netsuspend service upstart config file

description "netsuspend - a network-aware suspend utility"
author "craymen@gmail.com"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# netsuspend will fork and daemonize
expect fork

# Automatically restart process if crashed
respawn

# Start the process
exec netsuspend -D
